# trains.service
[Unit]
Description=Amtrak Train Timetable Server
After=network.target

[Service]
Type=simple
User=trains
WorkingDirectory=/home/trains/trains
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/trains/.local/bin/uv run server.py
Restart=always
RestartSec=5
RuntimeMaxSec=1h

[Install]
WantedBy=multi-user.target
